# overlay/etc/supervisor/conf.d/heartbeat.conf
[program:heartbeat]
command=/scripts/heartbeat.sh
stdout_events_enabled=true
stderr_events_enabled=true
